# tui/rpi_tui.service
[Unit]
Description=RPI TUI Application (runs on tty1)
After=network.target

[Service]
Type=simple
# Run as root (no User directive)
Environment=PYTHONUNBUFFERED=1
# Use venv python when installed to /opt/rpi-lab/.venv
ExecStart=/opt/rpi-lab/.venv/bin/python /opt/rpi-lab/tui/rpi_tui.py
Restart=on-failure
StandardInput=tty
StandardOutput=inherit
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes

[Install]
WantedBy=multi-user.target
