# gui/rpi_gui.service
[Unit]
Description=RPI Lab GUI Application (X11/Wayland)
After=graphical.target network-online.target
Wants=graphical.target

[Service]
Type=simple
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/mrmagic/.Xauthority"
Environment="PYTHONUNBUFFERED=1"
ExecStartPre=/bin/sleep 5
ExecStart=/opt/rpi-lab/.venv/bin/python /opt/rpi-lab/gui/rpi_gui.py
Restart=on-failure
RestartSec=10
User=mrmagic
WorkingDirectory=/opt/rpi-lab
# Allow I2C access for BME690 sensor without sudo
SupplementaryGroups=i2c

[Install]
WantedBy=graphical.target
