# sensors/mqtt_publisher.service
[Unit]
Description=BME690 MQTT Publisher for Home Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=mrmagic
WorkingDirectory=/opt/rpi-lab
ExecStart=/opt/rpi-lab/.venv/bin/python3 /opt/rpi-lab/sensors/mqtt_publisher.py
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Configuration file (edit /opt/rpi-lab/config/sensor.conf)
# See docs/HOME_ASSISTANT_MQTT.md for setup guide

# I2C access
SupplementaryGroups=i2c

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mqtt_publisher

[Install]
WantedBy=multi-user.target
